# dev_ops/stx-autospec-dashboard/autospec_checker.service
[Unit]
Description=personal

[Service]
Type=oneshot
ExecStart=/bin/python /usr/local/bin/main.py

[Install]
WantedBy=multi-user.target
