# dev_ops/stx-autospec-dashboard/web_report.service
[Unit]
Description=personal

[Service]
Type=oneshot
ExecStart=/bin/python /usr/local/bin/report.py

[Install]
WantedBy=multi-user.target
